# logitechmediaserver.conf
[supervisord]
nodaemon=true

[program:logitechmediaserver]
user=squeezeboxserver
command=/usr/sbin/squeezeboxserver --charset=utf8 --prefsdir /data/prefs --cachedir /data/cache --logdir /data/log
